# YAML/rules-emerging-threats/Privilege Escalation/defrag_deactivation.yaml
id: 4a6ce7fe-f2ca-55dd-a462-dd43a5212835
name: Defrag Deactivation
description: Detects the deactivation and disabling of the Scheduled defragmentation task as seen by Slingshot APT group
severity: Medium
queryFrequency: PT1H
queryPeriod: PT1H
triggerOperator: gt
triggerThreshold: 0
tactics:
- PrivilegeEscalation
- Execution
- Persistence
query: |-
  DeviceProcessEvents
  | where (ProcessCommandLine contains "/delete" or ProcessCommandLine contains "/change") and (ProcessCommandLine contains "/TN" and ProcessCommandLine contains "\\Microsoft\\Windows\\Defrag\\ScheduledDefrag") and FolderPath endswith "\\schtasks.exe"
version: 1.0.0
kind: Scheduled
relevantTechniques:
- T1053
entityMappings:
- entityType: Process
  fieldMappings:
  - identifier: CommandLine
    columnName: ProcessCommandLine
  - identifier: ProcessPath
    columnName: FolderPath
- entityType: File
  fieldMappings:
  - identifier: Directory
    columnName: FolderPath
requiredDataConnectors:
- connectorId: MicrosoftThreatProtection
  dataTypes:
  - DeviceProcessEvents
