# YAML/KQL/rules-threat-hunting/Execution/scheduled_task_created_filecreation.yaml
id: 440cea0c-b6dc-569d-9d88-ede7ee5d437e
name: Scheduled Task Created - FileCreation
description: Detects the creation of a scheduled task via file creation. - Normal behaviour on Windows
severity: Low
queryFrequency: PT1H
queryPeriod: PT1H
triggerOperator: gt
triggerThreshold: 0
tactics:
- Execution
- Persistence
- PrivilegeEscalation
query: |-
  DeviceFileEvents
  | where FolderPath contains ":\\Windows\\System32\\Tasks\\" or FolderPath contains ":\\Windows\\SysWOW64\\Tasks\\" or FolderPath contains ":\\Windows\\Tasks\\"
version: 1.0.0
kind: Scheduled
relevantTechniques:
- T1053
entityMappings:
- entityType: Process
  fieldMappings:
  - identifier: ProcessPath
    columnName: FolderPath
- entityType: File
  fieldMappings:
  - identifier: Directory
    columnName: FolderPath
requiredDataConnectors:
- connectorId: MicrosoftThreatProtection
  dataTypes:
  - DeviceFileEvents
