# YAML/rules/Privilege Escalation/scheduled_task_creation_with_curl_and_powershell_execution_combo.yaml
id: 9f065586-688e-5103-bcaf-67ca4f3d43c6
name: Scheduled Task Creation with Curl and PowerShell Execution Combo
description: Detects the creation of a scheduled task using schtasks.exe, potentially in combination with curl for downloading payloads and PowerShell for executing them. This facilitates executing malicious payloads or connecting with C&C server persistently without dropping the malware sample on the host. - Legitimate use of schtasks for administrative purposes. - Automation scripts combining curl and PowerShell in controlled environments.
severity: Medium
queryFrequency: PT1H
queryPeriod: PT1H
triggerOperator: gt
triggerThreshold: 0
tactics:
- PrivilegeEscalation
- Execution
- Persistence
- DefenseEvasion
- CommandAndControl
query: |-
  DeviceProcessEvents
  | where (ProcessCommandLine contains "curl " and ProcessCommandLine contains "http" and ProcessCommandLine contains "-o") and ((ProcessCommandLine contains " -create " or ProcessCommandLine contains " /create " or ProcessCommandLine contains " –create " or ProcessCommandLine contains " —create " or ProcessCommandLine contains " ―create ") and FolderPath endswith "\\schtasks.exe") and ProcessCommandLine contains "powershell"
version: 1.0.0
kind: Scheduled
relevantTechniques:
- T1053
- T1218
- T1105
entityMappings:
- entityType: Process
  fieldMappings:
  - identifier: CommandLine
    columnName: ProcessCommandLine
  - identifier: ProcessPath
    columnName: FolderPath
- entityType: File
  fieldMappings:
  - identifier: Directory
    columnName: FolderPath
requiredDataConnectors:
- connectorId: MicrosoftThreatProtection
  dataTypes:
  - DeviceProcessEvents
